# backend/distribution/linux/mpn-core_deb/DEBIAN/postinst
#!/bin/bash
set -e

echo "📦 Installing NAPS2 (safe apt install)..."

NAPS2_DEB="/usr/share/mpn-core/naps2-8.2.0-linux-x64.deb"

if [ -f "$NAPS2_DEB" ]; then
    apt-get update -y
    apt-get install -y "$NAPS2_DEB"
    echo "✅ NAPS2 installed successfully."
else
    echo "❌ ERROR: NAPS2 installer not found at $NAPS2_DEB"
    exit 1
fi

echo "⚠️ Installing main application (mpn-core)..."
systemctl daemon-reload
systemctl enable mpn-core.service
systemctl start mpn-core.service

update-desktop-database || true
gtk-update-icon-cache /usr/share/icons/hicolor || true

USER_DESKTOP="/home/$SUDO_USER/Desktop"
if [ -d "$USER_DESKTOP" ]; then
    cp /usr/share/applications/mpn-core.desktop "$USER_DESKTOP/"
    chmod +x "$USER_DESKTOP/mpn-core.desktop"
fi

echo "✅ Installation complete."
exit 0
