# backend/distribution/linux/mpn-core_deb/DEBIAN/prerm
#!/bin/bash
set -e
systemctl stop mpn-core.service || true
systemctl disable mpn-core.service || true
systemctl daemon-reload || true
exit 0
