# backend/distribution/linux/mpn/mpn-core_deb/DEBIAN/prerm
#!/bin/bash
set -e

echo "🛑 Stopping mpn-core..."

if systemctl is-active --quiet mpn-core.service; then
    systemctl stop mpn-core.service
fi
if systemctl is-enabled --quiet mpn-core.service; then
    systemctl disable mpn-core.service
fi

systemctl daemon-reload || true
exit 0
