# backend/build/com.naps2.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.naps2.service</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/naps2-service</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/naps2-service.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/naps2-service-error.log</string>
</dict>
</plist>
