# backend/distribution/mac/scripts/postinstall
#!/bin/bash

# -------------------------------
# Postinstall script for mpn-core
# -------------------------------

LOG_FILE="/tmp/mpn-postinstall.log"
BUILD_SCRIPT="/usr/local/mpn-core/build.sh"
SERVICE_NAME="mpn-core"

# Ensure the log exists and is writable
touch "$LOG_FILE"
chmod 666 "$LOG_FILE"

{
  echo "=== Postinstall started at $(/bin/date) ==="
  echo "Script running as user: $(whoami)"
  echo "Working directory: $(pwd)"
  echo "Environment PATH: $PATH"
  echo "-----------------------------------------"

  # -------------------------------
  # Run build.sh safely
  # -------------------------------
  if [ -f "$BUILD_SCRIPT" ]; then
      echo "→ Running build.sh..."
      /bin/chmod +x "$BUILD_SCRIPT"
      /bin/bash "$BUILD_SCRIPT" >> "$LOG_FILE" 2>&1 || echo "⚠️ build.sh failed" >> "$LOG_FILE"
  else
      echo "❌ build.sh not found at $BUILD_SCRIPT"
  fi

  echo "-----------------------------------------"
  echo "→ Attempting to start LaunchAgent for logged-in user (if any)..."

  # -------------------------------
  # Detect logged-in console user
  # -------------------------------
  CONSOLE_USER=$(stat -f "%Su" /dev/console)
  if [ -z "$CONSOLE_USER" ] || [ "$CONSOLE_USER" = "root" ]; then
      echo "⚠️ No logged-in GUI user detected. Service will start at next login."
  else
      USER_ID=$(id -u "$CONSOLE_USER")
      USER_HOME=$(eval echo "~$CONSOLE_USER")
      AGENT_PLIST="$USER_HOME/Library/LaunchAgents/com.$SERVICE_NAME.plist"

      echo "→ Detected GUI user: $CONSOLE_USER ($USER_ID)"
      echo "→ LaunchAgent path: $AGENT_PLIST"

      if [ -f "$AGENT_PLIST" ]; then
          echo "→ Bootstrapping LaunchAgent..."
          sudo -u "$CONSOLE_USER" launchctl bootout gui/$USER_ID "com.$SERVICE_NAME" 2>/dev/null || true
          sudo -u "$CONSOLE_USER" launchctl bootstrap gui/$USER_ID "$AGENT_PLIST" 2>/dev/null || true
          sudo -u "$CONSOLE_USER" launchctl kickstart -k gui/$USER_ID/com.$SERVICE_NAME 2>/dev/null || true
          echo "✅ LaunchAgent started for $CONSOLE_USER"
      else
          echo "⚠️ LaunchAgent plist not found for $CONSOLE_USER"
      fi
  fi

  echo "=== Postinstall finished at $(/bin/date) ==="
  echo "Logs written to: $LOG_FILE"
} >> "$LOG_FILE" 2>&1

# Always exit 0 so installer does not fail
exit 0
