# backend/distribution/linux/mpn-core-linux_deb/DEBIAN/postinst
#!/bin/bash
set -e

# Reload and enable service
systemctl daemon-reload
systemctl enable mpn-core-linux.service
systemctl start mpn-core-linux.service

# Update desktop database
if command -v update-desktop-database &>/dev/null; then
    update-desktop-database
fi

# Copy desktop shortcut to the installing user's Desktop
USER_DESKTOP="/home/$SUDO_USER/Desktop"
if [ -d "$USER_DESKTOP" ]; then
    cp /usr/share/applications/mpn-core-linux.desktop "$USER_DESKTOP/"
    chmod +x "$USER_DESKTOP/mpn-core-linux.desktop"
fi

exit 0
